# supervisord.conf
[supervisord]
user=root

nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=6
stopwaitsecs=2400
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
startretries=3
stopwaitsecs=2400
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.log





